# ops/templates/openclaw-restart-watch.timer
[Unit]
Description=Run OpenClaw restart-rate watcher every hour

[Timer]
OnCalendar=hourly
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
