# ops/templates/openclaw-backup.timer
[Unit]
Description=Run OpenClaw backup daily at 02:00 local

[Timer]
OnCalendar=*-*-* 02:00:00
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
