# ops/templates/openclaw-restart-watch.service.template
[Unit]
Description=OpenClaw restart-rate watcher
After=openclaw-gateway.service

[Service]
Type=oneshot
User=__OPENCLAW_USER__
Group=__OPENCLAW_GROUP__
EnvironmentFile=__OPENCLAW_ETC__/openclaw.env
ExecStart=__OPS_SCRIPT_DIR__/restart-watch.sh
